# scripts/profile.d/linux_armv7hf
WEBRTC_TARGET="libjingle_so"

GYP_DEFINES="\
OS=linux \
host_os=linux \
build_with_libjingle=1 \
build_with_chromium=0 \
clang=1 \
disable_iterator_debugging=1 \
disable_glibcxx_debug=1 \
target_arch=arm \
arm_version=7 \
arm_float_abi=hard \
arm_neon=1 \
arm_thumb=0 \
use_rtti=1 \
linux_fpic=1 \
libvpx_build_vp9=1 \
"
if [ "${WEBRTC_BUILD_TYPE}" == "Debug" ]; then
    GYP_DEFINES=$GYP_DEFINES" \
enable_tracing=1 \
"
fi

GYP_CROSSCOMPILE=1

DEB_ARCH=armhf

TOOLCHAIN=/usr/arm-linux-gnueabihf
