# provision/roles/chromium_dev/files/etc/profile.d/depot-tools.sh
#!/bin/sh
export PATH=/opt/depot_tools:${PATH}
